# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:laravel-worker]
command=php artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true